# Docker_Volume/Playbooks/Playbooks/ESXi_CIS/roles/cis/backup/files-bak/rhel8/config/aidecheck.service
[Unit]
Description=Aide Check

[Service]
Type=simple
SuccessExitStatus=7
ExecStart=/usr/sbin/aide --check

[Install]
WantedBy=multi-user.target
